# deployment/systemd/secrets-check.service
[Unit]
Description=Verify /run/secrets/ populated before Docker (ADR-045 Phase 6)
Documentation=https://github.com/achimdehnert/platform/blob/main/docs/adr/ADR-045-secrets-management.md
Before=docker.service
ConditionPathExists=/run/secrets

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'test -s /run/secrets/django_secret_key || (echo "FATAL: /run/secrets/django_secret_key missing — run CI/CD deploy first" >&2; exit 1)'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
