# infra/host-maintenance/infra-cleanup.timer
[Unit]
Description=Weekly host disk cleanup (prevention — keeps disk off 100%)

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
