# apps/bfagent/scripts/systemd/bfagent.service
[Unit]
Description=BF Agent Django Application
Documentation=https://github.com/yourusername/bfagent
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=notify

# User and Group
User=bfagent
Group=bfagent

# Working Directory
WorkingDirectory=/opt/bfagent

# Environment
Environment="PATH=/opt/bfagent/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DJANGO_ENV=production"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/bfagent/.env

# Main Process
ExecStart=/opt/bfagent/.venv/bin/gunicorn \
    config.wsgi:application \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --worker-class sync \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --access-logfile /var/log/bfagent/access.log \
    --error-logfile /var/log/bfagent/error.log \
    --log-level info \
    --capture-output

# Reload (for zero-downtime deploys)
ExecReload=/bin/kill -s HUP $MAINPID

# Restart Policy
Restart=always
RestartSec=10

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bfagent/media /opt/bfagent/staticfiles /opt/bfagent/logs

# Resource Limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bfagent

[Install]
WantedBy=multi-user.target
